# g733/99-g733-lights.rules
ACTION=="add", SUBSYSTEM=="hidraw", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="0ab5", RUN+="/usr/local/bin/g733-lights-off.sh"
